# dockerpress.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 10 * * * root /usr/local/bin/wp cache flush --path=/var/www/html
3 10 * * * root /usr/local/bin/mysql-optimize
